# slurm/all_up.service
[Unit]
Description=All up

[Service]
User=nimbix
ExecStart=/usr/local/bin/all_up.sh
